# FedoraScripts/archive/dotool.service
[Unit]
Description=dotool keyboard automation daemon

[Service]
ExecStart=/usr/local/bin/dotoold

[Install]
WantedBy=default.target
